# scripts/ai.fruityclaw.wq-decorrelate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.fruityclaw.wq-decorrelate</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/FruityClaw/.openclaw/workspace/worldquant-alpha-factory/scripts/run_decorrelate_pipeline.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/FruityClaw/.openclaw/workspace/worldquant-alpha-factory</string>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/FruityClaw/.openclaw/workspace/worldquant-alpha-factory/logs/launchd-decorrelate.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/FruityClaw/.openclaw/workspace/worldquant-alpha-factory/logs/launchd-decorrelate.err</string>
</dict>
</plist>
